# packaging/apk/post-install.sh
#!/bin/sh
rc-update add display-u6143 default
rc-service display-u6143 start
