# packaging/post-install.sh
#!/bin/bash
systemctl enable display-u6143.service
systemctl start display-u6143.service
